# clawsec-api.service
[Unit]
Description=ClawSec Cloud API
After=network.target

[Service]
Type=simple
User=openclaw
WorkingDirectory=/home/openclaw/clawsec-v2
ExecStart=/usr/bin/node api/src/server.js
Environment=PATH=/home/openclaw/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=CLAWSEC_PORT=3100
Environment=NODE_ENV=production
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target