# clawsec-sync.timer
[Unit]
Description=ClawSec Intel Sync Timer

[Timer]
OnCalendar=daily
OnBootSec=10min
Persistent=true

[Install]
WantedBy=timers.target